# deploy/digitalocean/systemd/celery-worker.service
[Unit]
Description=Celery Worker for NumerAI
After=network.target docker.service
Requires=docker.service
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=numerai
Group=numerai
WorkingDirectory=/opt/numerai
EnvironmentFile=/opt/numerai/.env.production

# Start the container
ExecStart=/usr/bin/docker-compose -f /opt/numerai/docker-compose.yml -f /opt/numerai/docker-compose.prod.yml run -d --rm --name numerai-celery-worker backend celery -A numerai worker -l info --concurrency=2

# Stop the container
ExecStop=/usr/bin/docker stop numerai-celery-worker || true
ExecStopPost=/usr/bin/docker rm numerai-celery-worker || true

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celery-worker

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
